# supervisord.conf
[supervisord]
nodaemon=true

[program:app]
command=node app.js

[program:nodejs-quickstart]
command=node /nodejs-quickstart/app.js
